# resources/udev-handler@.service
[Unit]
Description=Mount USB Drive on %i
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/etc/udev-handler/udev-handler.sh start %i
ExecStop=/etc/udev-handler/udev-handler.sh stop %i